# etc/rc.local
if [ -e /boot/setup]
then
  sudo apt-get install avahi-daemon
  sudo insserv avahi-daemon
  wget -q https://git.io/voEUQ -O /etc/avahi/services/multiple.service
  sudo /etc/init.d/avahi-daemon restart
fi
